# Payload.MegaDriveMini/opt/project_lunar/etc/udev/rules.d/90-networking.rules
ACTION=="add|remove", SUBSYSTEM=="net", RUN+="/etc/init.d/S97_PL_90_networking reload"